# systemd/fix-ports.service
[Unit]
Description=Fix forwarded torrent ports

[Service]
Type=oneshot
ExecStart=/bin/bash /home/<USER>/home-server/scripts/fix-ports.sh